# supervisord-nonroot.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid
logfile=/var/log/supervisor/supervisord.log
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:postgres]
command=/usr/lib/postgresql/9.5/bin/pg_ctl -D /var/lib/postgresql/9.5/main/ -o "-c config_file=/etc/postgresql/9.5/main/postgresql.conf" start
autostart=false
startsecs=0
autorestart=false
priority=1

[program:influxdb]
command=/usr/bin/influxd -config /etc/influxdb/influxdb.conf
priority=2

[program:pgwatch2]
command=/pgwatch2/pgwatch2
priority=3

[program:webpy]
command=python3 /pgwatch2/webpy/web.py
priority=4

[program:grafana]
command=/usr/sbin/grafana-server -homepath /usr/share/grafana/ --config=/etc/grafana/grafana.ini --pidfile=/var/run/grafana/grafana-server.pid cfg:default.paths.logs=/var/log/grafana cfg:default.paths.data=/var/lib/grafana cfg:default.paths.plugins=/var/lib/grafana/plugins
priority=5

[program:grafana_dashboard_setup]
command=/pgwatch2/bootstrap/set_up_grafana_dashboards.sh
priority=6
autorestart=false
